# docker-root/etc/services.d/tigervnc/run
#!/usr/bin/with-contenv bash
sudo -Hu user /usr/bin/vncserver :11 -fg -geometry $VNC_GEOMETRY -depth 24
