# docker-root/etc/services.d/websocketify/run
#!/usr/bin/with-contenv bash
sudo -Hu user /app/src/websockify/run 9001 127.0.0.1:5911